# brain.service
[Unit]
Description=Cecelia Intelligence Service (Semantic Search + Code Monitoring + Orchestration)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=xx
WorkingDirectory=/home/xx/dev/cecelia-semantic-brain
ExecStart=/usr/bin/docker-compose up
ExecStop=/usr/bin/docker-compose down
Restart=always
RestartSec=10

# Environment
Environment=COMPOSE_PROJECT_NAME=cecelia-semantic-brain

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cecelia-intelligence

[Install]
WantedBy=multi-user.target
